# systemd/wopr.service
[Unit]
Description=WOPR NeoPixel Pattern Manager
After=network.target

[Service]
Type=simple
User=achilles
WorkingDirectory=/home/achilles/WOPR/src
ExecStart=/home/achilles/WOPR/venv/bin/python3 /home/achilles/WOPR/src/run_service.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target